# rootfs/etc/services.d/budgetbakers_wallet/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting BudgetBakers Wallet app..."
exec python3 /app/src/main.py
